# synthesus_framework/aios/dist/macroqc-daemon.service
[Unit]
Description=MacroQC Quantum Simulation Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/macroqc/kernel/quantum_core.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
